# compose/development/app/supervisord.conf
[supervisord]
nodaemon=true

[program: node]
command=yarn start

[program:python]
command=python3 -m candis